# doc/contrib/init/dogecoind.conf
description "Corgicoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CORGCOIND_BIN="/usr/bin/corgicoind"
env CORGCOIND_USER="corgicoin"
env CORGCOIND_GROUP="corgicoin"
env CORGCOIND_PIDDIR="/var/run/corgicoind"
# upstart can't handle variables constructed with other variables
env CORGCOIND_PIDFILE="/var/run/corgicoind/corgicoind.pid"
env CORGCOIND_CONFIGFILE="/etc/corgicoin/corgicoin.conf"
env CORGCOIND_DATADIR="/var/lib/corgicoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # corgicoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CORGCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run corgicoind."
        echo "The setting must appear in $CORGCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Corgicoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CORGCOIND_PIDDIR"
    chmod 0755 "$CORGCOIND_PIDDIR"
    chown $CORGCOIND_USER:$CORGCOIND_GROUP "$CORGCOIND_PIDDIR"
    chown $CORGCOIND_USER:$CORGCOIND_GROUP "$CORGCOIND_CONFIGFILE"
    chmod 0660 "$CORGCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CORGCOIND_PIDFILE" \
    --chuid $CORGCOIND_USER:$CORGCOIND_GROUP \
    --exec "$CORGCOIND_BIN" \
    -- \
    -pid="$CORGCOIND_PIDFILE" \
    -conf="$CORGCOIND_CONFIGFILE" \
    -datadir="$CORGCOIND_DATADIR" \
    -disablewallet \
    -daemon

